# deploy/systemd/mmx-ops.timer
[Unit]
Description=Run MMX Ops Daily at 09:05

[Timer]
OnCalendar=*-*-* 09:05:00
Persistent=true

[Install]
WantedBy=timers.target
